# build/supervisord_autostart.conf
[program:xvfb]
priority=10
directory=/
command=sh /start_xvfb.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
environment=DISPLAY=":1"
stdout_logfile=/var/log/xvfb.log
redirect_stderr=true

[program:openbox]
priority=15
directory=/
command=/usr/bin/openbox
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
environment=DISPLAY=":1"
stdout_logfile=/var/log/openbox.log
redirect_stderr=true

[program:kde]
priority=15
directory=/
command=/usr/bin/plasmashell -n
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
environment=DISPLAY=":1"
stdout_logfile=/var/log/kde.log
redirect_stderr=true

[program:x11vnc]
priority=25
directory=/
command=sh -c "cat /root/.vncpasswd | xargs -i /usr/bin/x11vnc -listen 0.0.0.0 -rfbport 5900 -skip_lockkeys -shared -threads -noxdamage -repeat -passwd {} -display :1 -forever"
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/x11vnc.log
redirect_stderr=true

[program:novnc]
priority=25
directory=/noVNC-1.1.0/
command=/noVNC-1.1.0/utils/launch.sh --vnc 127.0.0.1:5900
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
stopasgroup=true

[program:nginx]
priority=25
directory=/
command=nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/nginx.log
redirect_stderr=true
stopasgroup=true

[program:web_autoresize]
priority=25
directory=/
command=/root/anaconda3/bin/python web_autoresize.py
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/web_autoresize.log
redirect_stderr=true
stopasgroup=true